# supervisor.conf
# Supervisor Configuration for Parking Payment System Queue Workers
# This configuration should be placed in /etc/supervisor/conf.d/parking-payment-system.conf

[program:parking-payment-system]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/parking-payment-system/artisan queue:work redis --sleep=3 --tries=3 --timeout=300
autostart=true
autorestart=true
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/parking-payment-system/queue-worker.log
stopwaitsecs=3600
environment=APP_ENV=production

# Restart the queue workers if they use too much memory
memmon_programs=parking-payment-system
memmon_eventlistener=TICK_60
